# discord_socket_client/root/etc/services.d/discord_socket_client/run
#!/usr/bin/with-contenv bash

exec \
  s6-setuidgid abc python3 /app/Legion-Chatbot/discord_socket_client/discord_socket_client.py